# systemd/iopaint.service
[Unit]
Description=IOPaint AI Image Processing Server
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/rus_bot
Environment="PATH=/home/YOUR_USERNAME/rus_bot/venv/bin"
ExecStart=/home/YOUR_USERNAME/rus_bot/venv/bin/iopaint start --model=lama --device=cpu --port=8085 --host=127.0.0.1
Restart=always
RestartSec=10

# Логирование
StandardOutput=append:/var/log/iopaint.log
StandardError=append:/var/log/iopaint-error.log

# Безопасность
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
